# protected/extensions/crontab/crontabs/my_crontab
* * * * * php /var/www/livescore/app/yiicmd.php test 'datetime'
0 0 * * * php /var/www/livescore/app/yiicmd.php test 'date'
